# service-units/windows/agentos-task.xml
<?xml version="1.0" encoding="UTF-8"?>
<!--
  AgentOS Windows Task Scheduler import — user-scope only.

  Install (PowerShell or cmd, as the invoking user, NO admin):
    schtasks /Create /XML service-units\windows\agentos-task.xml /TN AgentOS

  Uninstall:
    schtasks /Delete /TN AgentOS /F

  Safety contract — DO NOT REMOVE:
    The task invokes agentos.exe with the listen flag pinned to 127.0.0.1 so
    the gateway binds to loopback only. Network exposure is an explicit
    operator opt-in. To expose the gateway on the network, edit the
    Arguments element below to flip 127.0.0.1 to 0.0.0.0, delete and
    re-import the task with schtasks. Only run with the wildcard bind
    behind a trusted reverse proxy or VPN — agentos's first-class auth
    assumes loopback scope by default. The S19 boot WARNING
    (gateway.bind.public) will fire whenever the task is configured for
    wildcard bind.

    The Principal block uses GroupId=S-1-5-32-545 (BUILTIN\Users) and
    RunLevel=LeastPrivilege so the task runs in the invoking user's
    context with no admin elevation. Do NOT change these to SYSTEM,
    LocalService, or HighestAvailable — that would break the user-scope
    contract.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>AgentOS agent runtime gateway (loopback default)</Description>
    <URI>\AgentOS</URI>
  </RegistrationInfo>
  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <GroupId>S-1-5-32-545</GroupId>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Priority>7</Priority>
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>agentos.exe</Command>
      <Arguments>gateway run --listen 127.0.0.1</Arguments>
    </Exec>
  </Actions>
</Task>
